# ansible/roles/sessions/templates/kon-sync.service.j2
[Unit]
Description=Kon repo sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kon sync
User=root
Environment=GH_CONFIG_DIR=/opt/kon/.config/gh
StandardOutput=journal
StandardError=journal
